# obslight-fakeobs/packaging/debian.prerm
#!/bin/sh

%stop_on_removal fakeobs
if [ $1 -eq 0 ] ; then
/sbin/chkconfig --del fakeobs
fi
